# master/services/astromech-camera.service
[Unit]
Description=AstromechOS Camera Stream (mjpg-streamer)
After=network.target

[Service]
User=artoo
ExecStart=/usr/local/bin/astromech-camera-start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
