# master/services/r2d2-master.service
[Unit]
Description=R2-D2 Master Controller
After=network.target
Wants=network.target

[Service]
Type=simple
User=artoo
WorkingDirectory=/home/artoo/r2d2
ExecStart=/usr/bin/python3 /home/artoo/r2d2/master/main.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=r2d2-master

[Install]
WantedBy=multi-user.target
